# package/debian/mbridge.postrm
#!/usr/bin/env bash

systemctl stop mbridge
systemctl disable mbridge
systemctl daemon-reload
